# docker/rc.local
#!/bin/bash

if [ -f /root/.inside_docker_inside_qemu ]
    cd /root/repo
    git pull
    make inside
fi

exit 0
