# docker/dist/debian/rc.local
#!/bin/bash

test -x /root/runinside.sh && /root/runinside.sh &

exit 0
